# deployment/ansible/roles/kafka/templates/kafka.conf
start on runlevel [2345]
stop on runlevel [!2345]
respawn
exec {{ kafka_home }}/bin/kafka-server-start.sh {{ kafka_home }}/config/server.properties
